# crontab
*/5 * * * *  ~/monitor/cronwatch.sh
@reboot ~/monitor/monitor.py start > /dev/null
